# containers/profile.d/00-os.sh
#!/bin/sh

export SHELL="/bin/bash"
export LC_ALL="C.UTF-8"
export USER="$(whoami)"
export HOME="/home/${USER}"
